# com.tetsuzawa.dotfiles-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tetsuzawa.dotfiles-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>sync_claude.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__HOME__/dotfiles</string>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>StandardOutPath</key>
  <string>/tmp/dotfiles-sync.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/dotfiles-sync.stderr.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
